# .config/systemd/user/openrgb.service
[Unit]
Description=OpenRGB RGB Controller
After=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStart=/usr/local/bin/openrgb --server --profile blood-red
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
